# vufind/vufind-update.timer
[Unit]
Description=Update Vufind Sorl from remote API timer

[Timer]
OnCalendar=*-*-* *:15:00

[Install]
WantedBy=timers.target

